# apps/command-proxy/command-proxy.socket
[Socket]
ListenStream=127.0.0.1:1234
Accept=yes

[Install]
WantedBy=sockets.target
