# apps/command-proxy/command-proxy@.service
# TODO: Sandboxing

[Service]
DynamicUser=yes
ExecStart=/opt/sb/command-proxy
StandardInput=socket
StandardOutput=socket
StandardError=journal
